# .config/crontab
@reboot 0 5 * * * systemctl --user restart pipewire pipewire-pulse wireplumber
